# Flask_Gunicorn/supervisord.conf
[supervisord]
nodaemon=true

[program:gunicorn]
command=/usr/bin/gunicorn test:app -b 0.0.0.0:8080
directory=/usr/src/app
autostart=true
autorestart=true

