# Flask_Supervisord/supervisord.conf
[supervisord]
nodaemon=true

[program:app1]
command=python /usr/src/app1/app1.py
autostart=true
autorestart=true

[program:app2]
command=python /usr/src/app2/app2.py
autostart=true
autorestart=true

